# examples/wwff-db.timer
[Unit]
Description=Update WWFF cluster data periodically

[Timer]
OnUnitActiveSec=10m
OnBootSec=0s

[Install]
WantedBy=timers.target
